# debian/pifi.service
[Unit]
After=NetworkManager.service

[Service]
ExecStart=/usr/bin/pifi_ap_startup

[Install]
WantedBy=multi-user.target

